# scripts/polymarket-sync.timer
[Unit]
Description=Push paper trades CSV to GitHub every 15 minutes
Requires=polymarket-sync.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s

[Install]
WantedBy=timers.target
